# systemd/hannah-llama.service
[Unit]
Description=llama-server for Hannah (resident local LLM)
After=network.target

[Service]
Type=simple
# The launcher picks whichever model Hannah currently has selected.
ExecStart=/ssd/repos/z-web/time-poc/run_llama_server.sh
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
